# files/beegfs-client.service
[Service]
Type=forking
Restart=on-failure
RestartSec=5s
